# vm/etc/rc.local
cd /usr/tests
kyua test
kyua report
shutdown -p now
